# deploy/ticker-watchdog.service
[Unit]
Description=Edge ticker health watchdog

[Service]
Type=oneshot
ExecStart=/opt/edge-ticker/deploy/watchdog.sh
